# .config/systemd/user/backup-home.service
[Unit]
Description=Backup home directory

[Service]
ExecStart=/home/frank/.local/bin/backup-home

