# .config/systemd/user/sync-nextcloud.service
[Unit]
Description=Sync nextcloud with local machine

[Service]
ExecStart=/home/frank/.local/bin/sync-nextlcoud

