# .config/systemd/user/sync-nextcloud.service
[Unit]
Description=sync nextcloud with local machine

[Service]
ExecStart=/home/frank/.local/bin/sync-nextcloud

