# .config/systemd/user/sync-nextcloud.timer
[Unit]
Description=sync nextcloud with local machine

[Timer]
OnActiveSec=15min
OnUnitActiveSec=2w 

[Install]
WantedBy=timers.target

